# packaging/systemd/bsdm-cache-indexer.service
[Unit]
Description=BSDM-Proxy cache-indexer (Kafka to OpenSearch)
Documentation=https://github.com/onixus/bsdm-proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bsdm-proxy
Group=bsdm-proxy
EnvironmentFile=-/etc/bsdm-proxy/cache-indexer.env
ExecStart=/opt/bsdm-proxy/bin/cache-indexer
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
