# packaging/systemd/bsdm-proxy.service
[Unit]
Description=BSDM-Proxy HTTPS caching proxy
Documentation=https://github.com/onixus/bsdm-proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bsdm-proxy
Group=bsdm-proxy
EnvironmentFile=-/etc/bsdm-proxy/bsdm-proxy.env
ExecStart=/opt/bsdm-proxy/bin/proxy
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
